# supervisord.ini
[supervisord]

[program:internal]
command=python3 ./internal.py
redirect_stderr=True
stdout_logfile=./logs/internal.log

[program:main]
command=python3 ./main.py
redirect_stderr=True
stdout_logfile=./logs/main.log

[program:caddy]
command=bash ./caddy_run.sh
redirect_stderr=True
stdout_logfile=./logs/caddy.log
